# rootfs/etc/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:nginx]
command = /usr/sbin/nginx -c /etc/nginx/nginx.conf  -g 'daemon off;'
autostart = true
numprocs = 1
priority = 10
stdout_events_enabled = true
stderr_events_enabled = true
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:php-fpm]
command = /usr/local/sbin/php-fpm -F
autostart = true
numprocs = 1
stdout_events_enabled = true
stderr_events_enabled = true
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:sshd]
command = /usr/sbin/sshd -D -e
autostart = true
autorestart = true
startretries = 3
